# supervisord/templates/default/supervised-program.conf.erb
[program:<%= @program.name %>]
command = <%= @program.command %>
