# supervisord/templates/default/supervised-program.conf.erb
[program:<%= @program.name %>]
command = <%= @program.command %>
<% @program.each_setting do |name, setting| -%>
<%= name %> = <%= setting %>
<% end -%>
